# web/utils/moo-manual.service
[Service]
Type=simple
ExecStart=/opt/moo/backend

[Install]
WantedBy=multi-user.target
